# dns-lg.conf
# DNS-lg - start the DNS-lg
description "Start DNS-lg"
author  "SIDN Labs team"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn
respawn limit 10 5

chdir /home/miekg/dns-lg
exec /home/miekg/dns-lg/lg-dns -loc "NL, hosted at SIDN Labs, AS 1140" -mail "sidnlabs@sidn.nl" -res "Unbound with DNSSEC validation"
